# config/crontab
*/1 * * * * /bin/bash -l -c 'cd cloudexchange.org && app/updater.rb >>cron.log 2>&1'
